# fb-yv2-misc/misc-yv2-fb/service_files/xyz.openbmc_project.Chassis.Control.Misc.service
[Unit]
Description=Facebook Misc Ipmi Read
Wants=ipmb.service
After=ipmb.service

[Service]
Restart=always
RestartSec=3
ExecStart= echo nothing
Type=dbus
BusName=xyz.openbmc_project.Chassis.Event

[Install]
WantedBy=sysinit.target

